# opssat-package/nmf/nmf-package/control/prerm
#!/bin/sh
systemctl stop nmf_supervisor.service
systemctl disable nmf_supervisor.service
